# data/varlink-http-bridge.service.in
#  SPDX-License-Identifier: LGPL-2.1-or-later
#
#  This file is part of varlink-http-bridge.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.
#  SPDX-License-Identifier: LGPL-2.1-or-later

[Unit]
Description=Varlink HTTP Bridge
DefaultDependencies=no
Conflicts=shutdown.target

[Service]
Type=simple
ExecStart=@bindir@/varlink-http-bridge --bind 0.0.0.0:1031
Restart=on-failure
RestartSec=5

[Install]
WantedBy=network.target
