# data/varlink-httpd.service.in
#  SPDX-License-Identifier: LGPL-2.1-or-later
#
#  This file is part of varlink-http-bridge.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.
#  SPDX-License-Identifier: LGPL-2.1-or-later

[Unit]
Description=Varlink HTTP Bridge
DefaultDependencies=no
Conflicts=shutdown.target

[Service]
Type=simple
ExecStart=@bindir@/varlink-httpd --bind 0.0.0.0:1031
Restart=on-failure
RestartSec=5

# TLS credentials: well-known names from the credstore, renamed to the
# short names the service expects in $CREDENTIALS_DIRECTORY.
# Place certs in /etc/credstore/ and keys in /etc/credstore.encrypted/.
ImportCredential=varlink-httpd.tls.certificate:cert
ImportCredential=varlink-httpd.tls.key:key
ImportCredential=varlink-httpd.tls.trust:trust

# SSH authorized keys (see systemd.system-credentials(7))
ImportCredential=ssh.authorized_keys.root
ImportCredential=ssh.ephemeral-authorized_keys-all

[Install]
WantedBy=network.target
